# ffind-daemon.service
[Unit]
Description=Fast file finder daemon with real-time indexing
Documentation=man:ffind-daemon(8)
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ffind-daemon --foreground /home
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Uncomment and modify for your needs:
# ExecStart=/usr/local/bin/ffind-daemon --foreground --db /var/cache/ffind/index.db /home /var/www

# If using --db, create the cache directory:
# sudo mkdir -p /var/cache/ffind

[Install]
WantedBy=multi-user.target
